# systemd/grafana-server.service
[Unit]
Description=Grafana Server
Documentation=https://grafana.com/docs/
After=network.target influxdb.service
Wants=network.target

[Service]
Type=simple
User=grafana
Group=grafana
WorkingDirectory=/usr/share/grafana
ExecStart=/usr/sbin/grafana-server --config=/etc/grafana/grafana.ini --homepath=/usr/share/grafana
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=grafana

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/grafana /var/log/grafana

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment variables
Environment=GF_SECURITY_ADMIN_USER=admin
Environment=GF_SECURITY_ADMIN_PASSWORD=admin
Environment=GF_USERS_ALLOW_SIGN_UP=false

[Install]
WantedBy=multi-user.target
